# conf/systemd.service
[Unit]
Description=Liberchat encrypted chat application
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__
Environment=NODE_ENV=production
Environment=PORT=__PORT__
ExecStart=/opt/node_n/n/versions/node/18/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=inherit

# Sandboxing options to harden the service
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__INSTALL_DIR__ __DATA_DIR__

[Install]
WantedBy=multi-user.target